# etc/supervisor.d/supervisord-killswitch.ini
[program:killswitch]
command = /killswitch.sh
startsecs = 0
priority=30
autorestart = false
startretries = 1