# etc/supervisor.d/supervisord-stuck_check.ini
[program:stuck_check]
command = /stuck_check.sh
startsecs = 0
priority=40
numprocs=1
autostart=true
autorestart=true