# etc/supervisor.d/supervisord-openvpn.ini
[program:openvpn]
command=/start_vpn.sh
numprocs=1
priority=10
autostart=true
autorestart=true
